# NetGuard/services/analysis-tools-collector.service
[Unit]
Description=NetGuard Pro - Analysis Tools Collector Service (Master Orchestrator)
Documentation=https://github.com/netguard-pro
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/home/jarvis/NetGuard
ExecStart=/usr/bin/python3 /home/jarvis/NetGuard/scripts/analysis_tools_collector.py
Restart=always
RestartSec=10
StandardOutput=append:/home/jarvis/NetGuard/logs/system/analysis-tools-service.log
StandardError=append:/home/jarvis/NetGuard/logs/system/analysis-tools-service-error.log

# Allow sudo for analysis tools that need root
Environment="PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"

# Security settings
NoNewPrivileges=false
PrivateTmp=true

[Install]
WantedBy=multi-user.target

