# NetGuard/services/flask-dashboard.service
[Unit]
Description=NetGuard Pro - Flask Web Dashboard
Documentation=https://netguard.local/docs
After=network.target network-online.target
Wants=network-online.target

[Service]
Type=simple
User=jarvis
Group=jarvis
WorkingDirectory=/home/jarvis/NetGuard/web
ExecStart=/usr/bin/python3 /home/jarvis/NetGuard/web/app.py

# Restart configuration for reliability
Restart=always
RestartSec=10
StartLimitBurst=5
StartLimitIntervalSec=60

# Logging
StandardOutput=append:/home/jarvis/NetGuard/logs/system/flask-dashboard.log
StandardError=append:/home/jarvis/NetGuard/logs/system/flask-dashboard-error.log

# Security settings
NoNewPrivileges=true
PrivateTmp=true

# Resource limits
LimitNOFILE=10000
MemoryMax=500M
CPUQuota=50%

# Environment
Environment="PYTHONUNBUFFERED=1"
Environment="FLASK_ENV=production"

[Install]
WantedBy=multi-user.target

