# NetGuard/services/netsniff-collector.service
[Unit]
Description=NetGuard Pro - netsniff-ng High-Performance Collector Service
After=network-online.target
Wants=network-online.target

[Service]
User=root
Group=root
WorkingDirectory=/home/jarvis/NetGuard
ExecStart=/usr/bin/python3 /home/jarvis/NetGuard/scripts/netsniff_collector.py
Restart=always
RestartSec=10
StandardOutput=append:/home/jarvis/NetGuard/logs/system/netsniff-service.log
StandardError=append:/home/jarvis/NetGuard/logs/system/netsniff-service-error.log
LimitNOFILE=65535
MemoryMax=400M
CPUAccounting=true
CPUShares=384

[Install]
WantedBy=multi-user.target

