# NetGuard/services/network-capture.service
[Unit]
Description=NetGuard Pro - tcpdump Network Capture Service
Documentation=https://github.com/netguard-pro
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/home/jarvis/NetGuard
ExecStart=/home/jarvis/NetGuard/scripts/start_tcpdump.sh
Restart=always
RestartSec=10
StandardOutput=append:/home/jarvis/NetGuard/logs/system/tcpdump-service.log
StandardError=append:/home/jarvis/NetGuard/logs/system/tcpdump-service-error.log

# Security settings
NoNewPrivileges=false
PrivateTmp=true

[Install]
WantedBy=multi-user.target

