# NetGuard/services/tcpdump-collector.service
[Unit]
Description=NetGuard Pro - tcpdump Professional Packet Collector
Documentation=https://netguard.local/docs/tcpdump
After=network.target network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/home/jarvis/NetGuard

# Professional tcpdump collector with zero packet loss design
ExecStart=/usr/bin/python3 /home/jarvis/NetGuard/scripts/tcpdump_collector.py

# Restart configuration for reliability
Restart=always
RestartSec=10

# Resource limits for high-performance capture
LimitNOFILE=65536
Nice=-5

# Logging
StandardOutput=append:/home/jarvis/NetGuard/logs/system/tcpdump-service.log
StandardError=append:/home/jarvis/NetGuard/logs/system/tcpdump-service-error.log

# Security settings
NoNewPrivileges=false
PrivateTmp=true

[Install]
WantedBy=multi-user.target

