# .cursor-server/data/User/History/48e9646c/GU4R.service
[Unit]
Description=NetGuard Pro - ngrep Pattern Matching Collector Service
After=network-online.target
Wants=network-online.target

[Service]
User=root
Group=root
WorkingDirectory=/home/jarvis/NetGuard
ExecStart=/usr/bin/python3 /home/jarvis/NetGuard/scripts/ngrep_collector.py
Restart=always
RestartSec=10
StandardOutput=append:/home/jarvis/NetGuard/logs/system/ngrep-service.log
StandardError=append:/home/jarvis/NetGuard/logs/system/ngrep-service-error.log
LimitNOFILE=65535
MemoryMax=300M
CPUAccounting=true
CPUShares=256

[Install]
WantedBy=multi-user.target

