# .cursor-server/data/User/History/38f3dcaa/ThWB.service
[Unit]
Description=NetGuard Pro - Suricata EVE Log Collector Service
Documentation=https://github.com/netguard-pro
After=network.target suricata.service

[Service]
Type=simple
User=jarvis
WorkingDirectory=/home/jarvis/NetGuard
ExecStart=/usr/bin/python3 /home/jarvis/NetGuard/scripts/suricata_collector.py
Restart=always
RestartSec=10
StandardOutput=append:/home/jarvis/NetGuard/logs/system/suricata-collector-service.log
StandardError=append:/home/jarvis/NetGuard/logs/system/suricata-collector-service-error.log

# Security settings
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target

