# .cursor-server/data/User/History/568b069d/QnIL.service
[Unit]
Description=NetGuard Pro - Web Dashboard Service
Documentation=https://github.com/netguard-pro
After=network.target

[Service]
Type=simple
User=jarvis
WorkingDirectory=/home/jarvis/NetGuard/web
ExecStart=/usr/bin/python3 /home/jarvis/NetGuard/web/app.py
Restart=always
RestartSec=10
StandardOutput=append:/home/jarvis/NetGuard/logs/system/dashboard-service.log
StandardError=append:/home/jarvis/NetGuard/logs/system/dashboard-service-error.log

# Environment
Environment="FLASK_APP=app.py"
Environment="FLASK_ENV=production"

# Security settings
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target

