# .cursor-server/data/User/History/-74e493c9/mdhf.service
[Unit]
Description=NetGuard Pro - Enhanced tshark Collector Service
Documentation=https://github.com/netguard-pro
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/home/jarvis/NetGuard
ExecStart=/usr/bin/python3 /home/jarvis/NetGuard/scripts/tshark_collector.py
Restart=always
RestartSec=10
StandardOutput=append:/home/jarvis/NetGuard/logs/system/tshark-service.log
StandardError=append:/home/jarvis/NetGuard/logs/system/tshark-service-error.log

# Security settings
NoNewPrivileges=false
PrivateTmp=false

# Resource limits
LimitNOFILE=65535
CPUQuota=50%
MemoryMax=500M

# Environment
Environment="PYTHONUNBUFFERED=1"

[Install]
WantedBy=multi-user.target

